# systemd/dde-session-daemon.target.wants/dde-lock.timer
[Unit]
After=dde-session-daemon.target

[Timer]
OnActiveSec=5s
AccuracySec=100ms
RandomizedDelaySec=1s
RemainAfterElapse=no