# systemd/dde-session-initialized.target.wants/dde-session-daemon.target
[Unit]
Description=dde-session-daemon target
