# systemd/dde-session-initialized.target.wants/dde-session-daemon.target
[Unit]
Description=dde-session-daemon target

BindsTo=dde-session-daemon.service
After=dde-session-daemon.service
