# systemd/dde-session.target
[Unit]
Description=dde-session target
DefaultDependencies=no

BindsTo=graphical-session.target
Before=graphical-session.target

Requires=dde-session-manager.service
After=dde-session-manager.service

Requires=dde-session-initialized.target
After=dde-session-initialized.target
