# systemd/dde-session-manager.service.in
[Unit]
Description=dde-session-manager service
RefuseManualStart=yes
RefuseManualStop=yes
OnFailure=dde-session-shutdown.target
OnFailureJobMode=replace-irreversibly
CollectMode=inactive-or-failed

Requires=dde-session-manager.target
Before=dde-session-manager.target

[Service]
Type=notify
NotifyAccess=main
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/dde-session --systemd-service
ExecStopPost=@CMAKE_INSTALL_FULL_LIBEXECDIR@/dde-session-ctl --shutdown
Slice=services.slice
