# systemd/dde-session-pre.service
[Unit]
Description=dde pre init service
RefuseManualStart=yes
RefuseManualStop=yes

[Service]
Type=simple
ExecStart=/usr/bin/dde_update_dbus_env
