# systemd/services.slice
[Unit]
Description=DDE User Services Slice
Documentation=dde services

[Slice]
# resource control